# cron/etc/crontab
# updatedb every hour
0 * * * * root /sbin/updatedb
